# rtax57/host/rtax57-apply.service
[Unit]
Description=Re-apply RT-AX57 hardening and DNS policy
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=matthew
Environment=HOME=/home/matthew
ExecStart=/usr/local/bin/rtax57-apply.sh
